# bin/pointor.service
[Unit]
Description=Pointor Service

[Service]
Type=simple
ExecStart=/home/admin/service/pointor/bin/pointor.js
Restart=on-failure
RestartSec=10
User=root

[Install]
WantedBy=multi-user.target
